# scripts/com.lifelogger.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
每晚定时转写 — 23:00 自动跑流水线，结果写入 Vault

安装：
  bash /Users/zyongzhu/workbase/github/lifelogger/scripts/install_launchd.sh

查看日志：
  tail -f /tmp/lifelogger-pipeline.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lifelogger.pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/zyongzhu/workbase/github/lifelogger/.venv/bin/python</string>
        <string>/Users/zyongzhu/workbase/github/lifelogger/main.py</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/zyongzhu/workbase/github/lifelogger</string>

    <!-- 每天 23:00 运行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/lifelogger-pipeline.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/lifelogger-pipeline.log</string>
</dict>
</plist>
