# scripts/com.lifelogger.recorder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
录音守护进程 — 开机自启，持续录音

安装：
  bash /Users/zyongzhu/workbase/github/lifelogger/scripts/install_launchd.sh

卸载：
  launchctl unload ~/Library/LaunchAgents/com.lifelogger.recorder.plist

查看日志：
  tail -f /tmp/lifelogger-recorder.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lifelogger.recorder</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zyongzhu/workbase/github/lifelogger/scripts/run_recorder.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/zyongzhu/workbase/github/lifelogger</string>

    <!-- 开机自动启动 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 崩溃后自动重启 -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/lifelogger-recorder.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/lifelogger-recorder.log</string>
</dict>
</plist>
